# threatgate.service
[Unit]
Description=ThreatGate IOC Manager
After=network.target

[Service]
Type=simple
User=threatgate
Group=threatgate
WorkingDirectory=/opt/threatgate
ExecStart=/bin/bash /opt/threatgate/start.sh
Restart=always
RestartSec=5
Environment="PATH=/opt/threatgate/venv/bin:/usr/bin:/bin"

# Hardening (service runs as threatgate; data dir writable)
NoNewPrivileges=true
ReadWritePaths=/opt/threatgate/data
ProtectHome=read-only
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
