# threatgate-misp-sync.service
[Unit]
Description=ThreatGate MISP Sync - Pull IOCs from MISP
After=network.target

[Service]
Type=oneshot
User=threatgate
Group=threatgate
WorkingDirectory=/opt/threatgate
ExecStart=/opt/threatgate/venv/bin/python /opt/threatgate/misp_sync_job.py
Environment="PATH=/opt/threatgate/venv/bin:/usr/bin:/bin"

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/threatgate/data
ProtectHome=true
