# ziochub.service
[Unit]
Description=ZIoCHub IOC Manager
After=network.target

[Service]
Type=simple
User=ziochub
Group=ziochub
WorkingDirectory=/opt/ziochub
ExecStart=/bin/bash /opt/ziochub/start.sh
Restart=always
RestartSec=5
EnvironmentFile=-/opt/ziochub/data/ziochub.env
Environment="PATH=/opt/ziochub/venv/bin:/usr/bin:/bin"

# Hardening (service runs as ziochub; data dir writable)
NoNewPrivileges=true
ReadWritePaths=/opt/ziochub/data
ProtectHome=read-only
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
