# ziochub-redirect.service
[Unit]
Description=ZIoCHub HTTP to HTTPS Redirect
After=network.target ziochub.service
BindsTo=ziochub.service

[Service]
Type=simple
User=ziochub
Group=ziochub
WorkingDirectory=/opt/ziochub
ExecStart=/opt/ziochub/venv/bin/python /opt/ziochub/http_redirect.py
EnvironmentFile=-/opt/ziochub/data/ziochub.env
Environment="REDIRECT_HTTP_PORT=8080"
Restart=always
RestartSec=5

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true

[Install]
WantedBy=multi-user.target
